# supervisord.conf
[program:nginx]
command=/usr/local/nginx/sbin/nginx  -g "daemon off;"
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0

[program:php-fpm]
command=/usr/local/php/sbin/php-fpm --nodaemonize --allow-to-run-as-root
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0

[include]
files=/etc/supervisor/*.conf

[supervisord]
nodaemon=true

[supervisorctl]